# dist/rusty-sweeper-monitor.service
[Unit]
Description=Rusty Sweeper Disk Monitor
Documentation=man:rusty-sweeper(1)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/rusty-sweeper monitor
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

# Environment
Environment=DISPLAY=:0

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=default.target
